# modules/nastori_networks/files/wvdial.service
[Unit]
Description=wvdial service

[Service]
ExecStart=/usr/bin/wvdial
Restart=always

[Install]
WantedBy=multi-user.target
